# framework-led-daemon/systemd/framework-led-daemon.service
[Unit]
Description=Framework LED Matrix status daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/framework-led-daemon run
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
